# dot_config/systemd/user/media-shares-guard.service
[Unit]
Description=Auto-unmount Plex media shares when off the homelab LAN

[Service]
Type=oneshot
ExecStart=%h/.local/bin/scripts/media-shares guard
