# dot_config/systemd/user/playwright-mcp.service
[Unit]
Description=Playwright MCP Server
Documentation=https://github.com/microsoft/playwright-mcp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PLAYWRIGHT_BROWSERS_PATH=/opt/playwright-browsers
ExecStart=/usr/bin/distrobox enter playwright -- env PLAYWRIGHT_BROWSERS_PATH=/opt/playwright-browsers npx @playwright/mcp@latest --headless --port 8931
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=playwright-mcp

[Install]
WantedBy=default.target
