# dot_config/autostart/blueman.desktop
[Desktop Entry]
Type=Application
Name=Blueman
Comment=Bluetooth Manager
Icon=blueman
Exec=blueman-applet
Categories=System;
Hidden=true
